# init/supervisord.conf
[supervisord]
nodaemon=true

[program:rasadefault_server]
command=rasa run --cors "*" --port 5005

[program:rasa_actions_server]
command=rasa run actions --port 5055


